# resources/appimagelauncherd_openrc.in
#!/sbin/openrc-run

#OpenRC user service

description="AppImageLauncher daemon"
supervisor="supervise-daemon"
command="@CMAKE_INSTALL_PREFIX@/bin/appimagelauncherd"
