# resources/appimagelauncherfs.service.in
[Unit]
Description=AppImageLauncherFS daemon
Before=display-manager.service

[Service]
Type=exec
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/appimagelauncherfs
Restart=always
RestartSec=2

[Install]
# multi-user.target etc. seem to be buggy on some systems, it may not even exist on some systems
# can be checked with systemctl --user list-units --type target
# default seems to work better -- thanks @smac89 for the hint
WantedBy=default.target
