# resources/appimagelauncherfs.service.in
[Unit]
Description=AppImageLauncherFS daemon
Before=display-manager.service

[Service]
Type=exec
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/appimagelauncherfs
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
